# root/etc/services.d/nginx/run
#!/bin/sh
# SparkyFitness nginx s6 service

echo "[INFO] Starting nginx..."
exec nginx -g "daemon off;" -c /usr/local/etc/nginx/sparky.conf
